# deploy/guild-website-backup.timer
[Unit]
Description=Run the Shengshi website backup every day

[Timer]
OnCalendar=*-*-* 03:30:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
